# test_driven/deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for Django-TDD"

start on net-device-up
stop on shutdown

respawn

setuid elspeth
chdir /home/tsung/sites/djangotdd.com/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    --access-logfile ../access.log \
    --error-logfile ../error.log \
    superlists.wsgi:application
